# deploy/wsl2/systemd/napcat.service
[Unit]
Description=NapCat (Docker Compose)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
RestartSec=10
TimeoutStartSec=600

# 按你的实际路径修改（默认适配本仓库在 /root 下的布局）
Environment=NAPCAT_DIR=/root/napcat

Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

ExecStart=/root/bot/deploy/wsl2/bin/napcat-up.sh
ExecStop=/root/bot/deploy/wsl2/bin/napcat-down.sh

[Install]
WantedBy=multi-user.target
